# cron/crontab
*/15 6-9 * * 1-5 /app/gommutetime -from "$FROM" -to "$TO" >> "/app/data/to-${ID}.csv"
*/15 14-17 * * 1-5 /app/gommutetime -to "$FROM" -from "$TO" >> "/app/data/from-${ID}.csv"
